# auto_down.cron
*/1 * * * * /PycharmProjects/QBuyPro/auto_down.sh