# auto_down.cron
*/1 * * * * /usr/src/QBuyPro/auto_down.sh